# centos-7/jdk-8/etc/profile.d/java.sh
export JAVA_HOME=/opt/jdk-8
export PATH=$JAVA_HOME/bin:$PATH

